# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/stdout
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

; HTTP server. artisan serve is single-threaded per worker; PHP_CLI_SERVER_WORKERS
; (set in the Dockerfile) lets it handle concurrent requests.
[program:web]
command=php artisan serve --host=0.0.0.0 --port=8080
directory=/app
autostart=true
autorestart=true
startretries=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

; Queue worker (QUEUE_CONNECTION=database). Without this, dispatched jobs
; (broadcast emails, heartbeats) would never be processed.
[program:queue]
command=php artisan queue:work --sleep=3 --tries=3 --max-time=3600
directory=/app
autostart=true
autorestart=true
stopwaitsecs=3630
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=20

; Scheduler. schedule:work is the container-friendly equivalent of a
; "* * * * * schedule:run" cron entry; it stays in the foreground.
[program:scheduler]
command=php artisan schedule:work
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=30
